# config/.config/autostart/warpd.sh.desktop
[Desktop Entry]
Exec=/usr/local/bin/warpd
Icon=dialog-scripts
Name=warpd.sh
Path=
Type=Application
X-KDE-AutostartScript=true
