# config/autostart/unclutter.sh.desktop
[Desktop Entry]
Exec=/usr/bin/unclutter
Icon=dialog-scripts
Name=unclutter.sh
Path=
Type=Application
X-KDE-AutostartScript=true
